# Java_Rev3/dockerfiles/hdp_storm_node/conf/supervisor/stormworker-supervisord.conf
[program:storm-supervisor]
command=storm supervisor
directory=/home/storm
autorestart=true
user=storm

[program:storm-logviewer]
command=storm logviewer
autorestart=true
user=storm

